# hrmd.service
[Unit]
Description=Huygens Remote Manager daemon

[Service]
ExecStart=/etc/init.d/hrmd start

[Install]
WantedBy=multi-user.target

